# iam_user_sync.cron
# Synchronize IAM users with local system

SSH_AUTHORIZED_KEYS_DIR="@@SSH_AUTHORIZED_KEYS_DIR@@"
IAM_AUTHORIZED_GROUPS="@@IAM_AUTHORIZED_GROUPS@@"
LOCAL_GROUPS="@@LOCAL_GROUPS@@"
LOCAL_MARKER_GROUP="@@LOCAL_MARKER_GROUP@@"
PATH="@@PATH@@"

*/10 * * * * root @@INSTALL_PREFIX@@/bin/iam_user_sync
